# systemd/kis-supply-bot-build-candidates.service
[Unit]
Description=KIS Supply Bot Build Candidates

[Service]
Type=oneshot
WorkingDirectory=/opt/kis-supply-bot
EnvironmentFile=/opt/kis-supply-bot/.env
ExecStart=/opt/kis-supply-bot/kis-supply-bot build-candidates
User=root
