# systemd/kis-supply-bot-collect-1610.timer
[Unit]
Description=Run KIS Supply Bot Collect Daily at 16:10

[Timer]
OnCalendar=Mon..Fri 16:10:00
Persistent=true

[Install]
WantedBy=timers.target
