# systemd/kis-supply-bot-build-candidates.timer
[Unit]
Description=Run KIS Supply Bot Build Candidates at 18:30

[Timer]
OnCalendar=Mon..Fri 18:30:00
Persistent=true

[Install]
WantedBy=timers.target
